# onumchat.service
[Unit]
Description=OnumChat server
After=network.target

[Service]
Type=simple
User=tim
WorkingDirectory=/home/tim/onumchat
Environment=PORT=3010

# URL pública do site (links de aprovar/negar conta no e-mail)
Environment=APP_BASE_URL=https://onumchat.timgo.uk

# E-mail: preencha com seu Gmail e a senha de app (não use a senha normal)
# Senha de app: https://myaccount.google.com/apppasswords
Environment=SMTP_HOST=smtp.gmail.com
Environment=SMTP_PORT=587
Environment=SMTP_USER=SEU_EMAIL@gmail.com
Environment=SMTP_PASS=SUA_SENHA_DE_APP_AQUI
Environment=SMTP_FROM=SEU_EMAIL@gmail.com

ExecStart=/home/tim/.nvm/versions/node/v22.22.1/bin/node server.js
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
